# pack-scripts/boards/sinlinx-a33-lcd5/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
#ttyS2::respawn:/sbin/getty 115200 /dev/ttyS2
::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r

